# mgnt-console/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=gunicorn main:app -b 0.0.0.0:5000
directory=/var/app
environment=LC_ALL=C.UTF-8