# wyze-engine/conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:python]
command=/opt/conda/bin/python main.py
directory=/var/app
environment=LC_ALL=C.UTF-8