# pre-processor/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=pypy3 main.py
directory=/var/app
environment=LC_ALL=C.UTF-8