# wyze-router/conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:python]
command=pypy3 -m app
directory=/var/app
environment=LC_ALL=C.UTF-8