# wyze-router/conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:python]
command=pypy3 main.py
directory=/var/app
environment=LC_ALL=C.UTF-8