# pre-processor/supervisor-conf/python.conf
[program:python]
command=python3 demuxer.py
directory=/var/app
environment=LC_ALL=C.UTF-8
